# bin/daftbot.conf
env USER="www-data"
env SCRIPT_DIR="/var/www/private/daftbot"
env LOG_FILE="/var/log/daftbot.log"

description "Daftbot startup script"
author      "Tristan van Bokkem"

start on startup
stop on shutdown

script
    exec /usr/bin/sudo -u $USER /usr/bin/php $SCRIPT_DIR/daftbot.php >> $LOG_FILE
end script